# Linux_Codes/Services/communication_test.service
[Unit]
Description=Communication Service

[Service]
Restart=on-failure
ExecStart=/usr/bin/communication_test

[Install]
WantedBy=multi-user.target