# Linux_Codes/Services/HeartbeatMechanism_car.service
[Unit]
Description=Heartbeat Mechanism Car

[Service]
ExecStart=/usr/bin/heartbeatMechanism_test_car

[Install]
WantedBy=multi-user.target