# Linux_Codes/Services/MCU_Communication.service
[Unit]
Description=MCU Communication Service

[Service]
Restart=always
RestartSec=0.5sec
ExecStart=/usr/bin/MCU_Communication

[Install]
WantedBy=multi-user.target