# Linux_Codes/Services/MCU_Communication.service
[Unit]
Description=MCU Communication Service
Before=HeartbeatMechanism_car.service wifi_server.service lidar_communication.service light_sign_detection.service lane_detection.service linux_control.service

[Service]
Type=simple
ExecStart=/home/gtu-hazine/Desktop/AllRepos_15.09.2019/autonomous-vehicle/Linux_Codes/local_build/Output/MCU_Communication/MCU_Communication
Restart=always
RestartSec=1
User=gtu-hazine

[Install]
WantedBy=multi-user.target
