# Linux_Codes/Services/wifi_client.service
[Unit]
Description=Wifi Client

[Service]
ExecStart=/usr/bin/client_test
Restart=on-failure
RestartSec=1sec

[Install]
WantedBy=multi-user.target