# Linux_Codes/Services/wifi_client.service
[Unit]
Description=Wifi Client

[Service]
Restart=always
RestartSec=1sec
ExecStart=/usr/bin/client_test

[Install]
WantedBy=multi-user.target