# Linux_Codes/Services/zmq_server.service
[Unit]
Description=zeromq Server

[Service]
ExecStart=/usr/bin/server_test
Restart=on-failure
RestartSec=1sec

[Install]
WantedBy=multi-user.target